# services/strongbox.target
[Unit]
Description=Strongbox Decryption Target

[Install]
WantedBy=strongbox-mount.service
