# services/strongbox.target
[Unit]
Description=Strongbox Decryption Target
Before=strongbox-mount.service

[Install]
WantedBy=strongbox-mount.service
